# debian/postrm
#!/bin/bash
# postrm script for virtualgps

systemctl enable gpsd.socket
