# debian/postinst
#!/bin/bash
# postinst script for virtualgps

systemctl enable virtualgps.service
systemctl start virtualgps.service

