# debian/prerm
#!/bin/bash
# prerm script for virtualgps

systemctl stop virtualgps.service
systemctl disable virtualgps.service
